# build/supervisord/supervisord.conf
[supervisord]
nodaemon = true
environment = GOBIN="/usr/local/go/bin",GOPATH="/root/go",PATH="$PATH:$GOBIN:/root/go/bin"

[program:sshd]
command = /usr/sbin/sshd -D -f /etc/ssh/sshd_config